# person-detector/rootfs/etc/services.d/person-detector/finish
#!/command/with-contenv bashio
# ==============================================================================
# Take down the service
# ==============================================================================
bashio::log.info "Service stopped"
