# person-detector/rootfs/etc/services.d/person-detector/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the person detector service
# ==============================================================================

# Change to app directory
cd /app

bashio::log.info "Starting person detector..."

# Run the detector script
exec python3 /app/person_detector.py \
    --camera "$(bashio::config 'camera_entity')" \
    --confidence "$(bashio::config 'confidence_threshold')" \
    --interval "$(bashio::config 'scan_interval')"
