# services.d/person_detector/run
#!/usr/bin/with-contenv bash

# Set up logging
exec 1> >(s6-log n20 s2000000 /var/log/person_detector)
exec 2>&1

# Change to the application directory
cd /app

# Execute the Python script with the environment variables
exec python3 person_detector.py \
    --camera "${CAMERA_ENTITY}" \
    --confidence "${CONFIDENCE_THRESHOLD}" \
    --interval "${SCAN_INTERVAL}"